# scripts/utilities/crontab/apply-default-crontab.sh
#!/bin/bash
# apply crontab
cat /var/www/elmsln/scripts/server/crontab.txt >> /etc/crontab
